# tasks/LinuxBroadcast/deploy/usignals.conf
#!upstart
description "QCTF School 2016, usignals"

start on startup
stop on shutdown

chdir /home/broadcast
exec ./usignals.py

